# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:nginx]
command=nginx -g "daemon off;" -c /app/nginx.conf
autostart=true
autorestart=true
stdout_logfile=/tmp/nginx.log
stderr_logfile=/tmp/nginx.err

[program:ml_service]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/ml_service
environment=PYTHONPATH="/app/ml_service"
autostart=true
autorestart=true
stdout_logfile=/tmp/ml_service.log
stderr_logfile=/tmp/ml_service.err
startsecs=30

[program:backend]
command=npm start
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/tmp/backend.log
stderr_logfile=/tmp/backend.err
startsecs=10
